# etc/hotplug.d/iface/25-dnsmasq
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

export DNSMASQ_HOTPLUG=1
/etc/init.d/dnsmasq enabled && /etc/init.d/dnsmasq start
