# etc/rc.local
#!/bin/sh

chown -R www-data:www-data /etc/auth
chmod 0700 /etc/auth

tor
python /lib/update/create-crontab.py
crond

chown -R www-data:www-data /etc/TZ
chmod 0744 /etc/TZ

chown -R www-data:www-data /etc/last_update_check
chmod 0744 /etc/last_update_check


exit 0
