# etc/rc.local
#!/bin/sh

chmod 0700 /etc/auth

exit 0
