# etc/rc.local
#!/bin/sh

chmod 0700 /etc/auth

chown -R www-data:www-data /etc/auth /etc/openwireless-setup-state /etc/dropbear/authorized_keys
chmod o+rx /etc/dropbear


tor
mkdir -p /var/spool/cron/crontabs
echo '0 * * * * python /lib/update/update.py' > /tmp/crontab
crontab /tmp/crontab
crond


exit 0
