# etc/rc.local
#!/bin/sh

chmod 0700 /etc/auth

touch /etc/openwireless-setup-state
chown -R www-data:www-data /etc/auth /etc/openwireless-setup-state /etc/dropbear/authorized_keys
chmod o+rx /etc/dropbear

tor
python /lib/update/create-crontab.py
crond

exit 0
