# infra/etc/init/xvfb.conf
# Xvfb - Xvfb virtual X windows

description "Xvfb server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

exec /usr/bin/Xvfb :0 -screen 0 1024x768x24
